# conf_redhat/etc/profile.d/sentinella.sh
export PATH=$PATH:/usr/share/python/sentinella/bin
